# packaging/linux/scripts/postinst
#!/bin/sh
set -e
systemctl daemon-reload || true
systemctl enable pollen-dek || true
systemctl restart pollen-dek || true
