# packaging/linux/scripts/prerm
#!/bin/sh
set -e
systemctl stop pollen-dek || true
systemctl disable pollen-dek || true
